# root/etc/s6-overlay/s6-rc.d/jackett/run
#!/command/execlineb -P
s6-setuidgid daemon
/opt/Jackett/jackett -d ${XDG_CONFIG_HOME} --NoUpdates --NoRestart -x